# restserver/package/restserver/restserver.conf
start on runlevel [2345]
stop on runlevel [!2345]
respawn
respawn limit 30 90
exec /opt/restserver/restserver
#exec python /opt/package/restserver/run.py
